# avail_build/bootstrap/ansible/roles/avail_prerender/templates/avail_prerender.conf
[program:avail_prerender]
directory=/usr/local/avail/avail_prerender
command=node server.js
user=avail
stdout_logfile=/var/log/avail/avail_prerender.log
stderr_logfile=/var/log/avail/avail_prerender_error.log
autostart=true
autorestart=true
